# roles/web/files/cuckoo_super.conf
[program:cuckoo]
command = /home/cuckoo/cuckoo-1.2/web/gunicorn_start.sh
user = cuckoo
stdout_logfile = /var/log/supervisor/cuckoo.log
redirect_stderr = true
